# recipes-core/busybox/files/inittab
# /etc/inittab — usb-proxy appliance (BusyBox init)

# Bring the system up: mounts /proc, /sys, devtmpfs, populates volatile dirs
# (needed for the read-only rootfs) and starts mdev.
::sysinit:/etc/init.d/rcS

# The single purpose of this device. Respawn on exit covers device
# unplug/replug and crashes (replaces the old systemd Restart=always).
::respawn:/usr/bin/usb-proxy-run

# Clean shutdown / reboot.
::shutdown:/etc/init.d/rcK
::ctrlaltdel:/sbin/reboot

# A respawn getty on the serial console (ttyS0) is appended below this line
# automatically by the busybox recipe from SERIAL_CONSOLES, giving recovery
# access without a network.
